# etc/wigo.cron
* * * * * root if [ ! `pgrep wigo` ] ; then /etc/init.d/wigo start > /dev/null 2> /dev/null ; fi
